# sdcard/system/etc/init.d/getroot
#!/system/bin/sh

### BEGIN INIT INFO
# Exec: ready
# Type: permanent
### END INIT INFO

insmod /system/lib/modules/getroot.ko
